# roles/pulseaudio/templates/pasystray.service.j2
[Unit]
Description=pasystray (pulseaudio tray manager)

[Service]
ExecStart=/usr/bin/pasystray
Restart=on-failure
StartLimitInterval=10

[Install]
WantedBy=default.target
