# roles/xorg/templates/compton.service.j2
[Unit]
Description=Compton

[Service]
ExecStart=/usr/bin/compton
Restart=on-failure
StartLimitBurst=10

[Install]
WantedBy=default.target
