# roles/redshift/templates/redshift.service.j2
[Unit]
Description=Redshift

[Service]
ExecStart=/usr/bin/redshift
Restart=on-failure
StartLimitBurst=10

[Install]
WantedBy=default.target
