# config/systemd/clipster.service
[Unit]
Description=Clipster daemon

[Service]
ExecStart=/usr/bin/clipster -d
Restart=on-failure

[Install]
WantedBy=default.target
