# roles/dunst/files/dunst-resume.timer
[Unit]
Description=Resume notifications

[Timer]
OnCalendar=Sun *-*-* 13:00:00
Persistent=true

[Install]
WantedBy=timers.target
