# roles/dunst/files/dunst-pause.timer
[Unit]
Description=Pause notifications

[Timer]
OnCalendar=Sun *-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
