# roles/clipster/templates/clipster.service.j2
[Unit]
Description=Clipster

[Service]
ExecStart=/usr/bin/clipster -d
Restart=on-failure
StartLimitBurst=10

[Install]
WantedBy=default.target
