# roles/kdeconnect/templates/indicator-kdeconnect.service.j2
[Unit]
Description=Indicator-KDEConnect

[Service]
ExecStart=/usr/bin/indicator-kdeconnect
Restart=on-failure
StartLimitBurst=10

[Install]
WantedBy=default.target
